# roles/influxdb/templates/prometheus-influxdb-stats-exporter.service.j2
[Unit]
Description=InfluxDB Stats Exporter
After=network.target

[Service]
User=prometheus
ExecStart=/usr/local/bin/influxdb_stats_exporter

[Install]
WantedBy=multi-user.target
