# roles/pinger/templates/pinger.service.j2
[Unit]
Description=Pinger Service
After=network.target

[Service]
Type=simple
User=pinger
ExecStart=/usr/local/bin/pinger.sh

[Install]
WantedBy=multi-user.target
